# debian/univention-guacamole-config.postrm
#!/bin/sh

#DEBHELPER#
case "$1" in
remove)
	test -x /etc/init.d/slapd && invoke-rc.d slapd restart
	. /usr/share/univention-lib/base.sh
	stop_udm_cli_server
	;;

purge|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
	;;

*)
	echo "postrm called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

exit 0
